# root/etc/uci-defaults/luci-fluent
#!/bin/sh
# Register luci-theme-fluent in LuCI and set up theme config defaults

uci batch <<EOF
set luci.themes.fluent=/luci-static/fluent
commit luci
EOF

# Ensure the fluent config exists
if ! uci -q get fluent.@global[0] > /dev/null 2>&1; then
	uci batch << 'EOF'
	set fluent.global=global
	set fluent.global.mode='normal'
	set fluent.global.primary='#0078D4'
	set fluent.global.dark_primary='#4da6ff'
	set fluent.global.page_bg='#fafafa'
	set fluent.global.card_bg='#ffffff'
	set fluent.global.sidebar_bg='#ffffff'
	set fluent.global.dark_page_bg='#1b1b1b'
	set fluent.global.dark_card_bg='#2d2d2d'
	set fluent.global.dark_sidebar_bg='#1f1f1f'
	set fluent.global.progressbar_font='#2e2b60'
	set fluent.global.dark_progressbar_font='#d6d9e5'
	set fluent.global.font_weight='400'
	set fluent.global.font_size='14'
	set fluent.global.border_radius='4'
	set fluent.global.sidebar_width='260'
	set fluent.global.sidebar_style='default'
	set fluent.global.header_height='48'
	set fluent.global.login_bg='builtin'
	set fluent.global.login_blur='40'
	set fluent.global.login_transparency='0.92'
	set fluent.global.login_animation='1'
	set fluent.global.card_shadow='small'
	set fluent.global.card_border='1'
	set fluent.global.transition_speed='normal'
	set fluent.global.custom_select='1'
	set fluent.global.view_transition='1'
	set fluent.global.tab_animation='1'
	set fluent.global.prefers_reduced_motion='1'
	set fluent.global.loading_bar='1'
EOF
	uci commit fluent
fi

exit 0
